# packaging/debian/stacklab.service
[Unit]
Description=Stacklab
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=simple
User=stacklab
Group=stacklab
SupplementaryGroups=docker
EnvironmentFile=-/etc/stacklab/stacklab.env
WorkingDirectory=/srv/stacklab
ExecStart=/usr/lib/stacklab/bin/stacklab
Restart=on-failure
RestartSec=5
TimeoutStartSec=30
TimeoutStopSec=30
NoNewPrivileges=false
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/srv/stacklab /var/lib/stacklab

[Install]
WantedBy=multi-user.target
