# packaging/systemd/stacklab.service.example
[Unit]
Description=Stacklab
After=network-online.target docker.service
Wants=network-online.target docker.service

[Service]
Type=simple
User=stacklab
Group=stacklab
SupplementaryGroups=docker
EnvironmentFile=-/etc/stacklab/stacklab.env
WorkingDirectory=/opt/stacklab/app/current
ExecStart=/opt/stacklab/app/current/bin/stacklab
Restart=on-failure
RestartSec=5
TimeoutStartSec=30
TimeoutStopSec=30
# Must remain false to allow the opt-in Docker admin helper to elevate via sudo.
NoNewPrivileges=false
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
# Include /etc/docker when the opt-in Docker admin helper is enabled.
ReadWritePaths=/opt/stacklab /var/lib/stacklab /etc/docker

[Install]
WantedBy=multi-user.target
